# deploy/systemd/piguard-api.service
[Unit]
Description=PiGuard API Server (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/piguard
Environment=PYTHONUNBUFFERED=1
# Ensure source tree is importable if not installed as a package
Environment=PYTHONPATH=/opt/piguard/src
ExecStart=/opt/piguard/.venv/bin/python -m wids api --config /etc/piguard/wids.yaml
Restart=on-failure
RestartSec=2
# Run as root for simplicity; no added capabilities needed
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
